# tools/agent-mcp/config/secrets-expiry-check@.service
[Unit]
Description=Ratis secrets expiry check
After=network.target

[Service]
Type=oneshot
User=%i
ExecStart=/usr/local/bin/uv run --project /opt/ratis/tools/agent-mcp python /opt/ratis/tools/agent-mcp/scripts/check_expiry.py
EnvironmentFile=/etc/ratis/secrets-expiry.env
StandardOutput=journal
StandardError=journal
